# var/service/powertop/run
#!/bin/sh
exec powertop --auto-tune & pause
